# systemd/hoshinoyakata-allsky.service
[Unit]
Description=Hoshinoyakata AllSky AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/hoshinoyakata_agent
ExecStart=/home/pi/hoshinoyakata_agent/venv/bin/python3 /home/pi/hoshinoyakata_agent/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
