# systemd/hoshinoyakata-allsky.service
[Unit]
Description=Hoshinoyakata Allsky AI Agent
After=network.target

[Service]
WorkingDirectory=/home/pi/hoshinoyakata_allsky_agent_real_v2
ExecStart=/home/pi/hoshinoyakata_allsky_agent_real_v2/venv/bin/python /home/pi/hoshinoyakata_allsky_agent_real_v2/web/app.py
Restart=always
RestartSec=5
User=pi

[Install]
WantedBy=multi-user.target
